# supervisord-synapse.conf
[program:synapse]
command=/start-synapse.sh
numprocs=1
autostart=true
autorestart=true