# supervisord-immutant.conf
[program:immutant]
command=/start-immutant.sh
numprocs=1
autostart=true
autorestart=true
